# agent-manager/dockerfiles/start.conf
[program:sp-app]
command = sh /workspace/start.sh
autorestart = false

[program:sp-webterminal]
command = sh /workspace/webterminal.sh
autorestart = false